# files/asterisk.conf
[program:asterisk]
user=asterisk
command=/usr/sbin/asterisk
stopwaitsecs=30
